# native/memcached/carbonio-memcached.service
[Unit]
Description=Carbonio Memcached Daemon
PartOf=carbonio-proxy.target

[Service]
Slice=carbonio-proxy.slice
User=zextras
ExecStart=/opt/zextras/common/bin/memcached -U 0 -l 127.0.1.1,127.0.0.1 -p 11211 -o modern,drop_privileges
Restart=always

# Security hardening (strict profile - cache daemon)
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
NoNewPrivileges=yes
RemoveIPC=yes
# Required for dropping privileges and running as a different user
CapabilityBoundingSet=CAP_SETGID CAP_SETUID CAP_SYS_RESOURCE
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
MemoryDenyWriteExecute=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native
ProtectHostname=yes
ProtectClock=yes

[Install]
WantedBy=carbonio-proxy.target
