# native/postfix/carbonio-postfix.path
# SPDX-FileCopyrightText: 2025 Zextras <https://www.zextras.com>
# SPDX-License-Identifier: GPL-2.0-only

# Watch postfix configuration for changes made by configd

[Unit]
Description=Watch Carbonio Postfix config for changes
Documentation=https://zextras.com
After=carbonio-postfix.service
BindsTo=carbonio-postfix.service

[Path]
PathModified=/opt/zextras/common/conf/master.cf
Unit=carbonio-postfix-reload.service

[Install]
WantedBy=carbonio-mta.target
