# native/opendkim/carbonio-opendkim.service
[Unit]
Description=Carbonio OpenDKIM Daemon
After=carbonio-configd.service
Wants=carbonio-configd.service
PartOf=carbonio-mta.target

[Service]
Slice=carbonio-mta.slice
User=zextras
ExecStartPre=echo REWRITE opendkim | netcat -w 120 localhost 7171
ExecStart=/opt/zextras/common/sbin/opendkim -f -x /opt/zextras/conf/opendkim.conf -u zextras
# OpenDKIM uses SIGUSR1 for config reload (not SIGHUP)
ExecReload=/bin/kill -USR1 $MAINPID
RestartSec=3s
Restart=on-failure

# Security hardening (strict profile)
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
NoNewPrivileges=yes
CapabilityBoundingSet=
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
MemoryDenyWriteExecute=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native
ProtectHostname=yes
ProtectClock=yes

[Install]
WantedBy=carbonio-mta.target
